# etc/services.d/mosquitto/run
#!/usr/bin/execlineb -P
/usr/sbin/mosquitto -c /mosquitto/conf/mosquitto.conf
