# etc/services.d/crond/run
#!/usr/bin/execlineb -P
/usr/sbin/crond -f -l 8 -d 6
